# deploy/conductor-tick.service
[Unit]
Description=Clade Conductor Tick (Kamaji)
After=network.target

[Service]
Type=oneshot
User=ubuntu
ExecStart=/home/ubuntu/.config/clade/conductor-tick.sh
EnvironmentFile=/home/ubuntu/.config/clade/conductor.env
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600
