# deploy/conductor-tick.timer
[Unit]
Description=Run Conductor tick hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
