# deploy/ember.service
[Unit]
Description=Clade Ember Server (Oppy)
After=network.target

[Service]
Type=simple
User=ian
WorkingDirectory=/home/ian/projects/mol_diffusion/OMTRA_oppy
ExecStart=/home/ian/mambaforge/envs/clade/bin/clade-ember
Restart=always
RestartSec=5
Environment="EMBER_PORT=8100"
Environment="EMBER_BROTHER_NAME=oppy"
Environment="EMBER_WORKING_DIR=/home/ian/projects/mol_diffusion/OMTRA_oppy"
Environment="HEARTH_URL=https://44.195.96.130"
Environment="HEARTH_API_KEY=CHANGE_ME"
Environment="HEARTH_NAME=oppy"

[Install]
WantedBy=multi-user.target
